# keepalived/etc/openrc/keepalived
#!/sbin/openrc-run

command="/sbin/$SVCNAME"
command_args="$KEEPALIVED_OPTS"
pidfile="/var/run/$SVCNAME.pid"

depend() {
	need net
	after firewall
}

